# atlas200.service
[Unit]
Description=atlas200 

[Service]
Type=simple
ExecStart=/usr/local/bin/atlas200.sh

[Install]
WantedBy=multi-user.target
