# atlas200.service
[Unit]
Description=atlas200 

[Service]
Type=simple
ExecStart=/usr/local/HiAI/atlas200.sh
